# deploy/music-downloader.launchd.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.tal.music-downloader</string>
    <key>WorkingDirectory</key>
    <string>/Users/your-user/Documents/GitHub/Music Downloader</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>/Users/your-user/Documents/GitHub/Music Downloader/scripts/share-start.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/your-user/Documents/GitHub/Music Downloader/tmp/runtime/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/your-user/Documents/GitHub/Music Downloader/tmp/runtime/launchd.err.log</string>
  </dict>
</plist>
